# conf/supervisord/herbalife_celery.ini
[program:herbalife-celery]
directory=/opt/finkok/herbalife
environment=DJANGO_SETTINGS_MODULE='herbalife.settings.production'
command=/opt/.virtualenvs/herbalife/bin/python /opt/finkok/herbalife/manage.py celery -A herbalife worker --loglevel=info --settings=herbalife.settings.production --concurrency=10 --time-limit=300 --autoscale=24,12
stdout_logfile=/var/log/herbalife/celery/celery-herbalife.log
stderr_logfile=/var/log/herbalife/celery/celery-herbalife.err

autostart=false
autorestart=false
startsecs=5

priority=998
user=nginx

