# conf/supervisord/herbalife_sftpwatch.ini
[program:herbalife-sftpwatch]
directory=/opt/finkok/herbalife
environment=DJANGO_SETTINGS_MODULE='herbalife.settings.production',PATH='/opt/.virtualenvs/herbalife/bin'
command=python -u /opt/finkok/herbalife/watch.py -d /var/sftp/herbalife/ENTRADA
stdout_logfile=/var/log/herbalife/supervisor/watch-herbalife.log
stderr_logfile=/var/log/herbalife/supervisor/watch-herbalife.err
redirect_stderr=true
autostart=false
autorestart=false
startsecs=5
priority=998
user=herbalife


