# conf/supervisord/herbalife_beat.ini
[program:herbalife-beat]
directory=/opt/finkok/herbalife
environment=DJANGO_SETTINGS_MODULE='herbalife.settings.production'
command=/opt/.virtualenvs/herbalife/bin/python /opt/finkok/herbalife/manage.py celery -A herbalife beat --loglevel=info --settings=herbalife.settings.production
stdout_logfile=/var/log/herbalife/celery/celery-herbalife-beat.log
stderr_logfile=/var/log/herbalife/celery/celery-herbalife-beat.err

autostart=false
autorestart=false
startsecs=5

priority=998
user=nginx

